# profile.d/all.sh
LSOFT_BASE=/opt/linearsoft/toolbag

. ${LSOFT_BASE}/profile.d/bash-ls.sh
. ${LSOFT_BASE}/profile.d/bash-nano.sh
. ${LSOFT_BASE}/profile.d/putty-displayFix.sh
. ${LSOFT_BASE}/profile.d/x11-forward.sh