# oracle/systemd/orashutdown.service
[Unit]
Description=Oracle Service Shutdown catch-all
Before=oralsnr@__ORADEFDB__.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutSec=1800
User=__USER__
Group=__GROUP__
ExecStop=__BASE__/oracle/systemd/orashutdown.sh __ORADEFDB__

[Install]
WantedBy=multi-user.target